# rootfs/etc/profile.d/visual.sh
# Copyright 2016 Roman Egunkov <romanegunkov@gmail.com>

[ -z "$VISUAL" ] && VISUAL=vim
export VISUAL

[ -z "$EDITOR" ] && EDITOR=vim
export EDITOR
